# deploy/systemd/amp-dream.timer
[Unit]
Description=Run the AMP dream pass nightly (off-peak)

[Timer]
# Off-peak nightly. Persistent=true catches up if the machine was asleep.
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=amp-dream.service

[Install]
WantedBy=timers.target
